# mcp/thenewrich-mcp.service
[Unit]
Description=thenewrich MCP Server
After=network.target

[Service]
Type=simple
User=ubuntu
WorkingDirectory=/var/www/thenewrich/mcp
ExecStart=/usr/bin/node index.js
Restart=always
RestartSec=5
Environment=NODE_ENV=production

[Install]
WantedBy=multi-user.target
